# bat_check.service
[Unit]
Description=Runs the battery test script 

[Service]
ExecStart=/usr/bin/bat_test.sh
